# templates/berad.service.j2
[Unit]
Description=berad daemon
After=network.target

[Service]
Type=simple
User={{ beranode_user_name }}
ExecStart=berad start --home {{ beranode_rootdir }}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

[Service]
LimitNOFILE=65535